# config/etc/supervisord.conf
[program:visa-scheduler]
command=wf scheduler